# srs.service
[Unit]
Description=SRS Streaming Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/srs-stream

ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose restart

TimeoutStartSec=90
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=srs-stream

[Install]
WantedBy=multi-user.target
